# autostart/swapcaps.desktop
[Desktop Entry]
Type=Application
Name=swapcaps
Exec=~/swap.ch
StartupNotify=true
