# autostart/xcompmgr.desktop
[Desktop Entry]
Type=Application
Name=xcompmgr
Exec=xcompmgr
StartupNotify=true
